# crontab
# Fetch + Predict: every 15 minutes from 06:00 to 23:45
*/15 6-23 * * * /work/M6_Project_Anti_Food_Waste/shell/run_fetch_predictions.sh >> /work/M6_Project_Anti_Food_Waste/outputs/cron_fetch_predictions.log 2>&1

# Fetch + Predict: once at midnight (00:00)
0 0 * * * /work/M6_Project_Anti_Food_Waste/shell/run_fetch_predictions.sh >> /work/M6_Project_Anti_Food_Waste/outputs/cron_fetch_predictions.log 2>&1

# ML pipeline — nightly at 02:00 (retrain only triggers if PR AUC drops below threshold)
0 2 * * * /work/M6_Project_Anti_Food_Waste/shell/run_ml.sh >> /work/M6_Project_Anti_Food_Waste/outputs/cron_ml.log 2>&1
