# crontab
# Fetch + Predict: every 15 minutes from 06:00 to 23:45
*/15 6-23 * * * cd /app && python fetch_pipeline/run_fetch.py && python predict_pipeline/predict.py >> /app/outputs/cron_fetch_predictions.log 2>&1

# Fetch + Predict: once at midnight (00:00)
0 0 * * * cd /app && python fetch_pipeline/run_fetch.py && python predict_pipeline/predict.py >> /app/outputs/cron_fetch_predictions.log 2>&1

# ML pipeline — nightly at 02:00 (retrain only triggers if PR AUC drops below threshold)
0 2 * * * cd /app && python ml_pipeline/run_ml.py >> /app/outputs/cron_ml.log 2>&1
